# launchd/com.epstein-rag.app.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.epstein-rag.app</string>

  <key>WorkingDirectory</key>
  <string>/Users/macstudio/Documents/RAG/Epstein_Files_RAG_macstudio</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/macstudio/Documents/RAG/Epstein_Files_RAG_macstudio/scripts/run_native.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/macstudio/Documents/RAG/Epstein_Files_RAG_macstudio/runtime/app.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/macstudio/Documents/RAG/Epstein_Files_RAG_macstudio/runtime/app.err.log</string>
</dict>
</plist>
